# crontab
# Sample crontab entries
30 9 * * * packtpub_free.sh
